# deluge.conf
[program:deluge]
user = nobody
command = /home/nobody/start.sh