# deluge.conf
[program:deluge]
user = nobody
command = /usr/bin/start.sh