# deluge.conf
[program:deluged]
user = nobody
command=/usr/bin/deluged

[program:deluge-web]
user = nobody
command=/usr/bin/deluge-web