# build/deluge.conf
[program:deluge-script]
autorestart = false
startsecs = 0
user = nobody
command = /home/nobody/deluge.sh
umask = 000

[program:shutdown-script]
autorestart = false
startsecs = 0
user = root
command = /usr/local/bin/shutdown.sh '^deluge-web,^/usr/bin/python /usr/bin/deluged'
umask = 000
