# build/deluge.conf
[program:start-script]
autorestart = false
startsecs = 0
user = nobody
command = /home/nobody/start.sh
umask = 000
